# deploy/deb/postinst
#!/bin/bash
# =============================================================================
# The Stem - Post-installation Script (deb)
# =============================================================================
# Creates service user/group, sets capabilities, configures directories,
# and enables the systemd service.
# =============================================================================
set -e

SERVICE_USER="stem"
SERVICE_GROUP="stem"
BINARY_PATH="/usr/bin/stem"
CONFIG_DIR="/etc/stem"
DATA_DIR="/var/lib/stem"
LOG_DIR="/var/log/stem"

# Create service group if it doesn't exist
if ! getent group "$SERVICE_GROUP" > /dev/null 2>&1; then
    echo "Creating group: $SERVICE_GROUP"
    groupadd --system "$SERVICE_GROUP"
fi

# Create service user if it doesn't exist
if ! getent passwd "$SERVICE_USER" > /dev/null 2>&1; then
    echo "Creating user: $SERVICE_USER"
    useradd --system \
        --gid "$SERVICE_GROUP" \
        --home-dir "$DATA_DIR" \
        --no-create-home \
        --shell /usr/sbin/nologin \
        --comment "The Stem Network Performance Testing" \
        "$SERVICE_USER"
fi

# Create config directory
if [ ! -d "$CONFIG_DIR" ]; then
    mkdir -p "$CONFIG_DIR"
fi
chown root:"$SERVICE_GROUP" "$CONFIG_DIR"
chmod 750 "$CONFIG_DIR"

# Install default config if not present
if [ ! -f "$CONFIG_DIR/config.yaml" ]; then
    if [ -f /usr/share/stem/config.yaml ]; then
        cp /usr/share/stem/config.yaml "$CONFIG_DIR/config.yaml"
        chown root:"$SERVICE_GROUP" "$CONFIG_DIR/config.yaml"
        chmod 640 "$CONFIG_DIR/config.yaml"
    fi
fi

# Create environment file if not present
if [ ! -f "$CONFIG_DIR/environment" ]; then
    cat > "$CONFIG_DIR/environment" << 'EOF'
# Stem environment variables
# STEM_AUTH_USERNAME=admin
# STEM_AUTH_PASSWORD=changeme
# STEM_JWT_SECRET=generate-a-secure-random-string
# STEM_LICENSE_KEY=your-license-key
EOF
    chown root:"$SERVICE_GROUP" "$CONFIG_DIR/environment"
    chmod 600 "$CONFIG_DIR/environment"
fi

# Create data directory
if [ ! -d "$DATA_DIR" ]; then
    mkdir -p "$DATA_DIR"
fi
chown "$SERVICE_USER:$SERVICE_GROUP" "$DATA_DIR"
chmod 750 "$DATA_DIR"

# Create log directory
if [ ! -d "$LOG_DIR" ]; then
    mkdir -p "$LOG_DIR"
fi
chown "$SERVICE_USER:$SERVICE_GROUP" "$LOG_DIR"
chmod 750 "$LOG_DIR"

# Set capabilities for raw socket access (required for packet testing)
# - CAP_NET_RAW: Required for raw packet I/O, ICMP
# - CAP_NET_ADMIN: Required for interface control
# - CAP_NET_BIND_SERVICE: Required for binding to privileged ports
if command -v setcap &> /dev/null; then
    setcap 'cap_net_raw,cap_net_admin,cap_net_bind_service=+ep' "$BINARY_PATH" || {
        echo "Warning: Could not set capabilities on $BINARY_PATH"
        echo "The service may need to run as root for full functionality."
    }
else
    echo "Warning: setcap not found. Install libcap2-bin for non-root operation."
    echo "  sudo apt install libcap2-bin"
fi

# Configure firewall (ufw) if available and active
if command -v ufw >/dev/null 2>&1; then
    if ufw status 2>/dev/null | grep -q "Status: active"; then
        ufw allow 8080/tcp comment 'Stem WebUI HTTP' >/dev/null 2>&1 || true
        ufw allow 8443/tcp comment 'Stem WebUI HTTPS' >/dev/null 2>&1 || true
        echo "Firewall configured for Stem service (ports 8080, 8443)"
    fi
fi

# Reload systemd daemon
systemctl daemon-reload

# Enable service to start on boot
systemctl enable stem.service

# Start the service
if systemctl is-active --quiet stem.service; then
    echo "Restarting stem service..."
    systemctl restart stem.service
else
    echo "Starting stem service..."
    systemctl start stem.service || true
fi

echo ""
echo "=============================================="
echo "  The Stem installed successfully!"
echo "=============================================="
echo ""
echo "Web interface: http://localhost:8080"
echo ""
echo "Quick start:"
echo "  1. Edit /etc/stem/environment to set credentials"
echo "  2. Restart: sudo systemctl restart stem"
echo ""
echo "Commands:"
echo "  View logs:      journalctl -u stem -f"
echo "  Restart:        sudo systemctl restart stem"
echo "  Status:         sudo systemctl status stem"
echo "  Stop:           sudo systemctl stop stem"
echo "  CLI help:       stem --help"
echo ""

exit 0
