# deploy/deb/prerm
#!/bin/bash
# =============================================================================
# The Stem - Pre-removal Script (deb)
# =============================================================================
# Stops the service before package removal.
# =============================================================================
set -e

# Stop the service if running
if systemctl is-active --quiet stem.service 2>/dev/null; then
    echo "Stopping stem service..."
    systemctl stop stem.service || true
fi

# Disable the service
if systemctl is-enabled --quiet stem.service 2>/dev/null; then
    echo "Disabling stem service..."
    systemctl disable stem.service || true
fi

exit 0
